# cpu-temp-stable-path.rules
ACTION=="add|change", SUBSYSTEM=="hwmon", ATTR{name}=="coretemp", RUN+="/usr/local/bin/create-temp-symlinks"

